# rootfs/etc/services.d/rbfeeder/run
#!/usr/bin/with-contenv sh
#shellcheck shell=sh

if [ "$(uname -m)" = "x86_64" ]
then
    if [ "$VERBOSE_LOGGING" = "true" ]; then
        /usr/bin/qemu-arm-static /usr/bin/rbfeeder 2>&1 | \
            awk -W Interactive '{print "[rbfeeder] " $0}'
    else
        /usr/bin/qemu-arm-static /usr/bin/rbfeeder 2>&1 | \
            awk -W Interactive '! / Statistics updated every 30 seconds /' | \
            awk -W Interactive '! /Packets sent in the last 30 seconds: /' | \
            awk -W Interactive '! / Data sent: /' | \
            awk -W Interactive '! / Data received: /' | \
            awk -W Interactive '{print "[rbfeeder] " $0}'
    fi

elif [ "$(uname -m)" = "armv7l" ]
then
    if [ "$VERBOSE_LOGGING" = "true" ]; then
        /usr/bin/rbfeeder 2>&1 | \
            awk -W Interactive '{print "[rbfeeder] " $0}'
    else
        /usr/bin/rbfeeder 2>&1 | \
            awk -W Interactive '! / Statistics updated every 30 seconds /' | \
            awk -W Interactive '! /Packets sent in the last 30 seconds: /' | \
            awk -W Interactive '! / Data sent: /' | \
            awk -W Interactive '! / Data received: /' | \
            awk -W Interactive '{print "[rbfeeder] " $0}'
    fi

elif [ "$(uname -m)" = "armhf" ]
then
    if [ "$VERBOSE_LOGGING" = "true" ]; then
        /usr/bin/rbfeeder 2>&1 | \
            awk -W Interactive '{print "[rbfeeder] " $0}'
    else
        /usr/bin/rbfeeder 2>&1 | \
            awk -W Interactive '! / Statistics updated every 30 seconds /' | \
            awk -W Interactive '! /Packets sent in the last 30 seconds: /' | \
            awk -W Interactive '! / Data sent: /' | \
            awk -W Interactive '! / Data received: /' | \
            awk -W Interactive '{print "[rbfeeder] " $0}'
    fi

elif [ "$(uname -m)" = "aarch64" ]
then
    if [ "$VERBOSE_LOGGING" = "true" ]; then
        /usr/bin/rbfeeder 2>&1 | \
            awk -W Interactive '{print "[rbfeeder] " $0}'
    else
        /usr/bin/rbfeeder 2>&1 | \
            awk -W Interactive '! / Statistics updated every 30 seconds /' | \
            awk -W Interactive '! /Packets sent in the last 30 seconds: /' | \
            awk -W Interactive '! / Data sent: /' | \
            awk -W Interactive '! / Data received: /' | \
            awk -W Interactive '{print "[rbfeeder] " $0}'
    fi
else
    echo ""
    echo "ERROR!"
    echo "This build is running on an unsupported architecture ($(uname -m))."
    echo "Please raise an issue on this container's GitHub reporting this."
    echo ""
    sleep 60
    exit 1
fi